# vaultwarden-backup.timer
[Unit]
Description=Daily Vaultwarden backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
